# docker/s6-overlay/s6-rc.d/redis/run
#!/command/with-contenv sh
exec su-exec redis redis-server \
  --dir /data/redis --appendonly yes --bind 127.0.0.1 --save 60 1000
